# roles/de/files/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=conky
Exec=conky --daemonize
StartupNotify=false
Terminal=false