# roles/de/files/autostart/foot.desktop
[Desktop Entry]
Type=Application
Name=foot
Exec=foot --server
StartupNotify=false
Terminal=false