# roles/de/files/autostart/urxvtd.desktop
[Desktop Entry]
Type=Application
Name=urxvtd
Exec=urxvtd -q -f -o
StartupNotify=false
Terminal=false