# roles/packages/files/systemd/bitcoind.service
[Unit]
Description=Bitcoind from flatpak
After=network.target

[Service]
ExecStart=/usr/bin/flatpak run --filesystem="/run/tor/control.authcookie:ro" --command=bitcoind org.bitcoincore.bitcoin-qt
Restart=on-failure

[Install]
WantedBy=default.target