# roles/packages/files/systemd/flatpak-update.service
# https://github.com/flatpak/flatpak/issues/1399#issuecomment-403065567
[Unit]
Description=Flatpak update

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update -y
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7