# roles/packages/files/systemd/foot.service
[Unit]
Description=Foot server
Documentation=https://codeberg.org/dnkl/foot
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/foot --server

[Install]
WantedBy=sway-session.target