# roles/packages/files/systemd/ipfs.service
[Unit]
Description=IPFS
Wants=network-online.target
After=network-online.target

[Service]
Restart=always
User=ipfs
Group=ipfs
Environment=PODMAN_SYSTEMD_UNIT=%n
ReadWritePaths=/var/lib/ipfs
ExecStart=/usr/bin/podman start -a ipfs
ExecStop=/usr/bin/podman stop -t 1 ipfs
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
Nice=16
IOSchedulingPriority=5

[Install]
WantedBy=default.target