# roles/packages/files/systemd/jackett.service
[Unit]
Description=Jackett
Wants=network-online.target
After=network-online.target

[Service]
Restart=always
User=jackett
Group=jackett
Environment=PODMAN_SYSTEMD_UNIT=%n
ReadWritePaths=/var/lib/jackett /etc/jackett
ExecStart=/usr/bin/podman start -a jackett
ExecStop=/usr/bin/podman stop -t 1 jackett
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
Nice=16
IOSchedulingPriority=5

[Install]
WantedBy=default.target