# roles/packages/files/systemd/osmscout-server.service
[Unit]
Description=OSM Scout Server

[Service]
ExecStart=/usr/bin/flatpak run io.github.rinigus.OSMScoutServer --console
Restart=on-failure
Nice=15

[Install]
WantedBy=default.target