# roles/packages/files/systemd/pipx-upgrade.service
[Unit]
Description=Pipx upgrade

[Service]
Type=oneshot
ExecStart=/usr/bin/pipx upgrade-all