# roles/packages/files/systemd/pipx-upgrade.timer
[Unit]
Description=Pipx update

[Timer]
OnCalendar=*-*-* 6:00
RandomizedDelaySec=60m
Persistent=true

[Install]
WantedBy=timers.target