# roles/packages/files/systemd/pipx-upgrade.service
[Unit]
Description=Pipx upgrade

[Service]
Type=oneshot
ExecStart=pipx upgrade-all
AppArmorProfile=-pipx
Nice=19