# roles/packages/files/systemd/searx.service
[Unit]
Description=Searx
Wants=network-online.target
After=network-online.target

[Service]
Restart=always
User=searx
Group=searx
Environment=PODMAN_SYSTEMD_UNIT=%n
ReadWritePaths=/var/lib/searx
ExecStart=/usr/bin/podman start -a searx
ExecStop=/usr/bin/podman stop -t 1 searx
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
Nice=10

[Install]
WantedBy=default.target