# roles/packages/templates/systemd/XD.service.j2
[Unit]
Description=Standalone I2P BitTorrent Client
After=network.target

[Service]
ExecStart=/usr/local/bin/XD /home/{{ username }}/.config/XD/torrents.ini
Restart=on-failure

# https://gist.github.com/ageis/f5595e59b1cddb1513d1b425a323db04
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
DevicePolicy=closed
ProtectSystem=strict
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
IPAddressDeny=any
IPAddressAllow=localhost

[Install]
WantedBy=default.target
